# jobs/crontab
* * * * * /app/deployscripts/_run_metricsprecalculation.sh /app/master-thesis
* * * * * /app/deployscripts/_run_postcollector.sh /app/master-thesis